# backend/aria-yt-dlp-update.service
[Unit]
Description=Aria yt-dlp self-update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=chait
ExecStart=/home/chait/MusicAppIOS/backend/update-yt-dlp.sh
# Allow the script's `systemctl restart aria-backend` without an interactive
# password by adding a sudoers rule, e.g.:
#   chait ALL=(root) NOPASSWD: /usr/bin/systemctl restart aria-backend
