# backend/aria-healthcheck.service
[Unit]
Description=Aria backend health probe
After=network-online.target

[Service]
Type=oneshot
User=chait
# Set ARIA_ALERT_WEBHOOK to a Slack/Discord/ntfy webhook to get pinged on
# failure. ARIA_HEALTH_URL defaults to the local instance.
# Environment=ARIA_ALERT_WEBHOOK=https://hooks.slack.com/services/XXX
ExecStart=/home/chait/MusicAppIOS/backend/healthcheck.sh
